# templates/rails-restart.path
[Path]
PathModified={{ rails__root_dir }}/revisions.log

[Install]
WantedBy=multi-user.target
