# templates/rails.service
[Unit]
After=network.target
Description={{ rails__general_descr }}

[Service]
ExecStart=/bin/true
ExecReload=/bin/true
RemainAfterExit=yes
Type=oneshot
WorkingDirectory={{ rails__root_dir }}/current

[Install]
WantedBy=multi-user.target
